# config/cron.hourly/sync_cpan
#!/bin/bash

/usr/bin/rsync -az --delete --no-motd rsync://ftp.funet.fi/CPAN /var/www/kineticode.com/cpan

