# config/cron.hourly/sync_cpan
#!/bin/bash

/usr/bin/rsync -az --delete --no-motd rsync://cpan-rsync.perl.org/CPAN/ /var/www/kineticode.com/cpan

