# config/cron.hourly/sync_pgxn
#!/bin/bash

/usr/bin/rsync -az --delete rsync://master.pgxn.org/pgxn /var/www/kineticode.com/pgxn

