# config/cron.daily/backupldap
#!/bin/sh

cd /var/backups
/usr/sbin/slapcat -b dc=kineticode,dc=com -l kineticode.ldif
/usr/bin/savelog -qjc 7 kineticode.ldif
