# config/cron.daily/backupimap
#!/bin/bash

# Must have /root/.offlineimaprc
/usr/bin/offlineimap
