# config/cron.daily/backuppostgres
#!/bin/bash

cd /var/backups
/usr/local/pgsql/bin/pg_dumpall -U postgres > postgresql.sql
/usr/bin/savelog -qjc 7 postgresql.sql
