# config/cron.daily/backupsvn
#!/bin/sh

cd /var/backups
/usr/local/bin/hot-backup.py /var/svn . > /dev/null
DUMP=`ls -d svn-*`
mv $DUMP svn
tar jcf svn.tar.bz2 svn
/usr/bin/savelog -qlc 7 svn.tar.bz2
rm -rf svn
