# config/cron.daily/backupsvn
#!/bin/sh

cd /var/backups
/usr/bin/hot-backup.py /usr/local/svnroot .
