# config/cron.daily/syncjsan
#!/bin/bash

#rsync -avz --exclude=/js/script/mirror.js rsync://master.openjsan.org/jsan /var/www/back.openjsan.org
